# systemd/embedded-monthly-builds-daily.service
[Unit]
Description=Embedded Monthly Builds Daily Cycle
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=nayab
WorkingDirectory=/home/nayab/embedded-monthly-builds/embedded-monthly-builds
Environment=HOME=/home/nayab
ExecStart=/usr/bin/python3 /home/nayab/embedded-monthly-builds/embedded-monthly-builds/scripts/run_daily_cycle.py --write-state

[Install]
WantedBy=multi-user.target
